# cron.d/perch
# Run perch tick every 2 hours
0 */2 * * * /app/run-perch.sh
