# cron.d/perch
# Run perch tick every 5 minutes
*/5 * * * * /app/scripts/run-perch.sh
